# desktop/gatepath-netns-helper/data/gatepath-netns-helper.service
[Unit]
Description=Gatepath captive-portal netns helper
Documentation=https://github.com/bearyjd/gatepath
# Need NetworkManager to be up before we can query "is this captive" — that
# query lands in 5b.4. Keeping the dependency declared from the start so
# socket activation orders correctly the moment we add it.
After=NetworkManager.service dbus.service
Requires=dbus.service

[Service]
Type=dbus
BusName=cc.grepon.Gatepath.NetNsHelper
ExecStart=/usr/libexec/gatepath-netns-helper
# Fixed PATH so the helper's `ip`/`iw`/`wpa_supplicant`/DHCP-client/`kill`
# lookups can't be redirected by an inherited or attacker-influenced PATH in
# a root process. systemd already supplies a default PATH (not inherited), but
# pinning it here is defence-in-depth for the binaries the helper now execs.
Environment=PATH=/usr/sbin:/usr/bin:/sbin:/bin
# Helper needs CAP_NET_ADMIN to run `ip netns ...` and move interfaces; it
# runs as root so that root's CAP_NET_ADMIN is implicit. We could drop to
# a dedicated user with AmbientCapabilities=CAP_NET_ADMIN, but distros that
# package this can make that decision per their security posture.
User=root

# Hardening — these are tight by default; loosen only if shown necessary.
NoNewPrivileges=false
# (Cannot set NoNewPrivileges=true because we exec `ip` from iproute2.
#  iproute2 doesn't need new privileges, but the test of suid bits across
#  distros has produced false positives. Reconsider once 5b switches to
#  direct netlink.)
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# AF_PACKET is required by the DESK-002 in-netns connectivity bring-up: DHCP
# clients (dhclient/udhcpc) use raw AF_PACKET sockets to send/receive DHCP
# before an address exists, and wpa_supplicant uses AF_PACKET for EAPOL.
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6 AF_PACKET
# NOTE (DESK-002): the helper now spawns network-using CHILDREN in the gatepath
# netns — wpa_supplicant, the DHCP client, and the portal WebView runner. Those
# children share this unit's cgroup, so a cgroup-level `IPAddressDeny=any` would
# block the captive sign-in itself (DHCP broadcast + the portal page). Network
# confinement is therefore the *netns*'s job (the children only ever see the
# captive Wi-Fi link), not the cgroup's. `IPAddressDeny=any` is intentionally
# NOT set for this reason; re-evaluate with `IPAddressAllow` scoping during the
# BLOCKER-DESK-003 hardware validation.
LockPersonality=true
# NOTE (DESK-003): MemoryDenyWriteExecute is left OFF because the portal WebView
# child (WebKitGTK) JITs JavaScript and needs W+X pages. Confirm on hardware
# whether a per-child relaxation is preferable to disabling it unit-wide.
MemoryDenyWriteExecute=false
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# State directory holds the audit log (helper-audit.jsonl). Mode 0750 so
# the audit log is reachable by root and the gatepath group only — closes
# a temporal-pattern leak to other local users that the default 0755 +
# default umask 0644 would expose. The audit log file itself is opened
# 0640 (see src/audit_log.rs).
StateDirectory=gatepath
StateDirectoryMode=0750

[Install]
# D-Bus activation handles startup; no Alias / WantedBy needed unless an
# operator wants to pre-start. The .service file in /usr/share/dbus-1/system-services/
# tells the bus which systemd unit to invoke.
